# service/user/.config/systemd/user/pon-arduino-cloud-agent.service
[Unit]
Description=Arduino Cloud Agent

[Service]
ExecStart=%h/ArduinoCloudAgent/Arduino_Cloud_Agent

[Install]
WantedBy=default.target
